# roles/satellite-receptor/templates/receptor@.service.j2
[Unit]
Description=Receptor Node for %i
After=network.target

[Service]
ExecStart=/usr/bin/receptor --node-id %i -c {{receptor_config_dir}}/%i/receptor.conf -d {{receptor_data_dir}}/%i node
